# packaging/deb/DEBIAN/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl stop migasfree-agent
    systemctl disable migasfree-agent
fi

exit 0
